# ops/alumni-agent.timer
[Unit]
Description=Run Nova SBE alumni networker agent every 15 minutes
Requires=alumni-agent.service

[Timer]
# Fire 2 minutes after boot, then every 15 minutes.
OnBootSec=2min
OnUnitActiveSec=15min
# If the system was off, run once shortly after boot to catch up.
Persistent=true
Unit=alumni-agent.service

[Install]
WantedBy=timers.target
